# services/dpmpv2.service
[Unit]
Description=DPMP v2 (non-docker)
After=network-online.target
Wants=network-online.target

[Service]
WorkingDirectory=%h/dpmp
Environment=DPMP_CONFIG=%h/dpmp/dpmp/config_v2.json
ExecStart=%h/dpmp/.venv/bin/python -u %h/dpmp/dpmp/dpmpv2.py
Restart=always
RestartSec=2
StandardOutput=append:%h/dpmp/dpmpv2_run.log
StandardError=append:%h/dpmp/dpmpv2_run.log

[Install]
WantedBy=default.target
